# config/tradingbot-monitor.service
[Unit]
Description=Trading Bot Monitor Service
After=network.target tradingbot.service

[Service]
Type=simple
User=tradingbot
Group=tradingbot
WorkingDirectory=/opt/tradingbot
Environment=PYTHONPATH=/opt/tradingbot
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/tradingbot/venv/bin/python monitor_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/tradingbot/logs/monitor.log
StandardError=append:/opt/tradingbot/logs/monitor.error.log

[Install]
WantedBy=multi-user.target 